# conf/mosdns.ini
[program:mosdns]
directory=/opt/dns
command=./scripts/start-dns.sh
autostart=yes
startsecs=5