# conf/smartdns.ini
[program:smartdns]
directory=/opt/dns
command=./smartdns -c smartdns.conf -f -x
autostart=yes
startsecs=2